# specs/cricket/cricket.cron
MAILTO=root
*/5 * * * * cricket @appdir@/collect-subtrees -cf=@etcdir@/subtree-sets normal #>/dev/null 2>&1
*/5 * * * * apache find @cachedir@ -type f -mmin +10 -print0 | xargs -0 -r rm

